# ops/crontab
# supercronic schedule for the VPS-hosted heartbeat (docker compose --profile heartbeat up -d)
# 00:00 UTC = 08:00 SGT, buffer inside the 08:30 promise
7 0 * * * cd /app && /opt/venv/bin/python -m agent.morning
